# scripts/99-nanoclock.rules
# Runs early at boot... (onboard *might* be mounted at that point)
KERNEL=="loop0", RUN+="/usr/local/NanoClock/bin/nanoclock-launcher.sh"
